# templates/default/services/upstart/eventserver.conf.erb
description "PIO EventServer"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [!2345]

kill timeout 20
setuid <%= @user %>
setgid <%= @group %>

limit nofile <%= @nofile %> <%= @nofile %>

env JAVA_HOME=<%= @java_home %>
env PIO_LOG_DIR<%= @logdir %>

script
    # source default daemon settings
    . <%= @environment_file %>

    exec <%= @piodir %>/bin/pio eventserver --port $PIO_EVENTSERVER_PORT
end script
