# templates/default/services/upstart/hadoop.conf.erb
description "Hadoop Services wrapper task"
start on filesystem and net-device-up IFACE!=lo
stop on starting rc RUNLEVEL=[!2345]

console none

pre-start script
    for s in hdfs-namenode hdfs-datanode; do
        start $s || :
    done
end script

post-stop script
    for s in hdfs-datanode hdfs-namenode; do
        stop $s || :
    done
end script
