# templates/default/services/upstart/hdfs-datanode.conf.erb
description "Hadoop DataNode Service"

start on started hdfs-namenode
stop on stopped hdfs-namenode

kill timeout 20
setuid <%= @user %>
setgid <%= @group %>

limit nofile <%= @nofile %> <%= @nofile %>

env JAVA_HOME=<%= @java_home %>
exec <%= @appdir %>/bin/hdfs datanode
