# templates/default/services/upstart/spark-worker.conf.erb
description "Spark Worker Service"

start on started spark-master
stop on stopped spark-master

kill timeout 20
setuid aml
setgid hadoop

limit nofile <%= @nofile %> <%= @nofile %>

env JAVA_HOME=<%= @java_home %>
env SPARK_HOME=<%= @home_prefix %>/spark
env DAEMON_ARGS=

exec $SPARK_HOME/sbin/start-spark.sh worker $DAEMON_ARGS
