# templates/default/services/systemd/eventserver.service.erb
[Unit]
Description=PIO EventServer
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
Environment=PIO_EVENTSERVER_PORT=7070
EnvironmentFile=<%= @initenv %>
WorkingDirectory=/var/log/eventserver
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @home_prefix %>/pio/bin/pio eventserver --port $PIO_EVENTSERVER_PORT

[Install]
WantedBy=multi-user.target
