# templates/default/services/systemd/eventserver.service.erb
[Unit]
Description=PIO EventServer
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
EnvironmentFile=-<%= @environment_file %>
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @piodir %>/bin/pio eventserver --port $PIO_EVENTSERVER_PORT

[Install]
WantedBy=multi-user.target
