# templates/default/services/systemd/eventserver.service.erb
[Unit]
Description=PIO EventServer
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
Environment=EVENTSERVER_PORT=7070
Environment=PIO_LOG_DIR=/var/log/eventserver
EnvironmentFile=-<%= @environment_file %>
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @appdir %>/bin/pio eventserver --port $EVENTSERVER_PORT

[Install]
WantedBy=multi-user.target
