# templates/default/services/systemd/hbase-master.service.erb
[Unit]
Description=HBase Master Server
After=syslog.target network.target
PartOf=hbase.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
Environment=LD_LIBRARY_PATH=<%= @hadoopdir %>/lib/native
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @appdir %>/bin/hbase master start

[Install]
WantedBy=multi-user.target
