# templates/default/services/systemd/hbase.service.erb
[Unit]
Description=HBase Services Launcher
After=syslog.target network.target
Requires=hbase-master.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
