# templates/default/services/systemd/hbase-master.service.erb
[Unit]
Description=HBase Master Server
After=syslog.target network.target
PartOf=hbase.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @home_prefix %>/hbase/bin/hbase master start

[Install]
WantedBy=multi-user.target
