# templates/default/services/systemd/hdfs-namenode.service.erb
[Unit]
Description=Hadoop NameNode Service
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
ExecStart=/usr/local/hadoop/bin/hdfs namenode

[Install]
WantedBy=multi-user.target
