# templates/default/services/systemd/hdfs-namenode.service.erb
[Unit]
Description=Hadoop NameNode Service
After=syslog.target network.target
PartOf=hadoop.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
ExecStart=<%= @home_prefix %>/hadoop/bin/hdfs namenode
ExecStartPre=/bin/sh -c "[ -f <%= @version_file %> ] || <%= @home_prefix %>/hadoop/bin/hdfs namenode -format"

[Install]
WantedBy=multi-user.target
