# templates/default/services/systemd/spark.service.erb
[Unit]
Description=Spark Services Launcher
After=syslog.target network.target
Requires=spark-master.service spark-worker.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
