# templates/default/services/systemd/spark-worker.service.erb
[Unit]
Description=Spark Worker Service
After=syslog.target network.target
Requires=spark-master.service
PartOf=spark.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=JAVA_HOME=<%= @java_home %>
Environment=SPARK_HOME=<%= @home_prefix %>/spark
Environment=DAEMON_ARGS=
LimitNOFILE=<%= @nofile %>
ExecStart=<%= @home_prefix %>/spark/sbin/start-spark.sh worker $DAEMON_ARGS

[Install]
WantedBy=multi-user.target
